# break-reminder.service
[Unit]
Description=Break reminder

[Service]
ExecStart=/usr/bin/notify-send --hint int:transient:1 "Take a break"
