# break-reminder.service
[Unit]
Description=Break reminder

[Service]
ExecStart=/usr/bin/notify-send "Take a break"
